# deployment/multi-mixed-scanner.service
[Unit]
Description=Multi-Asset Mixed Scanner (BTC, ETH, EUR/USD - 15m/1h/4h)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading-bot
Environment="PATH=/home/ubuntu/trading-bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ubuntu/trading-bot/venv/bin/python3 /home/ubuntu/trading-bot/main_multi_symbol.py --config /home/ubuntu/trading-bot/config/multi_mixed.json
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/trading-bot/logs/multi_mixed.log
StandardError=append:/home/ubuntu/trading-bot/logs/multi_mixed.log

[Install]
WantedBy=multi-user.target
