# deployment/multi-crypto-scalp-scanner.service
[Unit]
Description=Multi-Crypto Scalping Scanner (BTC, ETH, SOL - 1m/5m/15m)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading-bot
Environment="PATH=/home/ubuntu/trading-bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ubuntu/trading-bot/venv/bin/python3 /home/ubuntu/trading-bot/main_multi_symbol.py --config /home/ubuntu/trading-bot/config/multi_crypto_scalp.json
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/trading-bot/logs/multi_crypto_scalp.log
StandardError=append:/home/ubuntu/trading-bot/logs/multi_crypto_scalp.log

[Install]
WantedBy=multi-user.target
