# deployment/us100-scanner.service
[Unit]
Description=US100/NASDAQ Trading Scanner (5m/15m/1h timeframes)
Documentation=file:///home/predeshen/telegramscalperbot/DEPLOYMENT_GUIDE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/predeshen/telegramscalperbot
ExecStart=/usr/bin/python3 /home/predeshen/telegramscalperbot/main_us100.py
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=append:/home/predeshen/telegramscalperbot/logs/us100_scanner.log
StandardError=append:/home/predeshen/telegramscalperbot/logs/us100_scanner.log

# Environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONDONTWRITEBYTECODE=1"
EnvironmentFile=/home/predeshen/telegramscalperbot/.env

# Security hardening
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=no
ProtectHome=no
ReadWritePaths=/home/predeshen/telegramscalperbot/logs /home/predeshen/telegramscalperbot/config

# Resource limits
LimitNOFILE=4096
MemoryMax=512M
CPUQuota=50%

# Timeout settings
TimeoutStartSec=30
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
