# deployment/us30-momentum-scanner.service
[Unit]
Description=US30 Momentum Trading Scanner
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=predeshen
Group=predeshen
WorkingDirectory=/home/predeshen/telegramscalperbot
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /home/predeshen/telegramscalperbot/main_us30.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=us30-momentum-scanner

# Resource limits
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
