# deployment/us30-momentum-scanner.service
[Unit]
Description=US30 Momentum Trading Scanner
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading-bot
Environment="PATH=/home/ubuntu/trading-bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ubuntu/trading-bot/venv/bin/python3 /home/ubuntu/trading-bot/main_us30.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/trading-bot/logs/us30_scanner.log
StandardError=append:/home/ubuntu/trading-bot/logs/us30_scanner.log

[Install]
WantedBy=multi-user.target
