# deployment/us30-scalp-scanner.service
[Unit]
Description=US30 Scalping Scanner (5m/15m timeframes)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/telegramscalperbot
ExecStart=/usr/bin/python3 /home/ubuntu/telegramscalperbot/us30_scanner/main_us30_scalp.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/telegramscalperbot/logs/us30_scalp_scanner.log
StandardError=append:/home/ubuntu/telegramscalperbot/logs/us30_scalp_scanner.log

# Environment variables
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/ubuntu/telegramscalperbot/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=4096
MemoryMax=512M

[Install]
WantedBy=multi-user.target
